# deploy/systemd/usb-encryption-enforcer-ui.service
[Unit]
Description=USB Encryption Enforcer UI bridge (notifications)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/lib/usb-encryption-enforcer/.venv/bin/python3 -u /usr/libexec/usb-encryption-enforcer-ui
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
