# deploy/systemd/usb-enforcer-ui.service
[Unit]
Description=USB Encryption Enforcer UI bridge (notifications)
After=graphical-session.target dbus.service
PartOf=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/libexec/usb-enforcer-ui
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
