# debian-bundled/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create Python virtual environment
        if [ ! -d /usr/lib/usb-encryption-enforcer/.venv ]; then
            echo "Creating Python virtual environment..."
            python3 -m venv /usr/lib/usb-encryption-enforcer/.venv
        fi
        
        # Upgrade pip
        /usr/lib/usb-encryption-enforcer/.venv/bin/pip install --upgrade pip
        
        # Install Python dependencies from bundled wheels
        echo "Installing Python dependencies from bundled wheels..."
        /usr/lib/usb-encryption-enforcer/.venv/bin/pip install \
            --no-index \
            --find-links=/usr/lib/usb-encryption-enforcer/wheels \
            pyudev pydbus typing-extensions
        
        # Reload systemd and udev
        systemctl daemon-reload || true
        udevadm control --reload-rules || true
        udevadm trigger || true
        
        # Enable and start daemon
        systemctl enable usb-encryption-enforcerd.service || true
        systemctl start usb-encryption-enforcerd.service || true
        
        # Enable user service for all users
        mkdir -p /etc/systemd/user/default.target.wants
        ln -sf /usr/lib/systemd/user/usb-encryption-enforcer-ui.service \
            /etc/systemd/user/default.target.wants/usb-encryption-enforcer-ui.service || true
        
        echo "USB Encryption Enforcer (bundled) has been installed and started."
        echo "Configuration file: /etc/usb-encryption-enforcer/config.toml"
        ;;
esac

#DEBHELPER#

exit 0
