# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create Python virtual environment with system site packages
        # (allows access to system python3-gi, python3-gi-cairo, etc.)
        if [ ! -d /usr/lib/usb-enforcer/.venv ]; then
            echo "Creating Python virtual environment..."
            python3 -m venv --system-site-packages /usr/lib/usb-enforcer/.venv
        fi
        
        # Upgrade pip
        /usr/lib/usb-enforcer/.venv/bin/pip install --upgrade pip
        
        # Install Python dependencies from PyPI
        echo "Installing Python dependencies from PyPI..."
        /usr/lib/usb-enforcer/.venv/bin/pip install \
            pyudev>=0.24.0 \
            pydbus>=0.6.0 \
            typing-extensions>=4.8.0 \
            python-magic>=0.4.27 \
            pdfplumber>=0.10.0 \
            python-docx>=1.0.0 \
            openpyxl>=3.1.0 \
            python-pptx>=0.6.0 \
            odfpy>=1.4.0 \
            py7zr>=0.20.0 \
            rarfile>=4.1 \
            fusepy>=3.0.1 \
            xlrd>=2.0.0 \
            olefile>=0.46 \
            extract-msg>=0.41.0 \
            striprtf>=0.0.26 \
            xlrd>=2.0.0 \
            olefile>=0.46 \
            extract-msg>=0.41.0 \
            striprtf>=0.0.26
        
        # Migrate config file (add new options if missing)
        if [ -f /etc/usb-enforcer/config.toml ]; then
            if ! grep -q "allow_plaintext_write_with_scanning" /etc/usb-enforcer/config.toml; then
                sed -i '/allow_luks1_readonly = true/a allow_plaintext_write_with_scanning = false  # Allow write to unencrypted drives if content scanning is enabled' \
                    /etc/usb-enforcer/config.toml || true
            fi
        fi
        
        # Reload systemd and udev
        systemctl daemon-reload || true
        udevadm control --reload-rules || true
        udevadm trigger || true
        
        # Enable and start daemon
        systemctl enable usb-enforcerd.service || true
        systemctl start usb-enforcerd.service || true
        
        # Enable user service for all users (graphical-session.target)
        mkdir -p /etc/systemd/user/graphical-session.target.wants
        ln -sf /usr/lib/systemd/user/usb-enforcer-ui.service \
            /etc/systemd/user/graphical-session.target.wants/usb-enforcer-ui.service || true

        # Refresh desktop entries and icon cache for menu integration
        if command -v update-desktop-database >/dev/null 2>&1; then
            update-desktop-database /usr/share/applications >/dev/null 2>&1 || true
        fi
        if command -v gtk-update-icon-cache >/dev/null 2>&1; then
            gtk-update-icon-cache -f -t /usr/share/icons/hicolor >/dev/null 2>&1 || true
        fi
        
        # Reload user systemd daemon for all logged-in users
        for uid in $(loginctl list-users --no-legend | awk '{print $1}'); do
            systemctl --user --machine="${uid}@" daemon-reload 2>/dev/null || true
        done
        
        echo "USB Encryption Enforcer has been installed and started."
        echo "Configuration file: /etc/usb-enforcer/config.toml"
        echo ""
        echo "Note: Desktop notifications require the usb-enforcer-ui service to be running."
        echo "Run 'systemctl --user start usb-enforcer-ui' to start it now, or log out and back in."
        ;;
esac

#DEBHELPER#

exit 0
