# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create Python virtual environment
        if [ ! -d /usr/lib/usb-enforcer/.venv ]; then
            echo "Creating Python virtual environment..."
            python3 -m venv /usr/lib/usb-enforcer/.venv
        fi
        
        # Upgrade pip
        /usr/lib/usb-enforcer/.venv/bin/pip install --upgrade pip
        
        # Install Python dependencies from PyPI
        echo "Installing Python dependencies from PyPI..."
        /usr/lib/usb-enforcer/.venv/bin/pip install \
            pyudev>=0.24.0 \
            pydbus>=0.6.0 \
            typing-extensions>=4.8.0
        
        # Reload systemd and udev
        systemctl daemon-reload || true
        udevadm control --reload-rules || true
        udevadm trigger || true
        
        # Enable and start daemon
        systemctl enable usb-enforcerd.service || true
        systemctl start usb-enforcerd.service || true
        
        # Enable user service for all users (graphical-session.target)
        mkdir -p /etc/systemd/user/graphical-session.target.wants
        ln -sf /usr/lib/systemd/user/usb-enforcer-ui.service \
            /etc/systemd/user/graphical-session.target.wants/usb-enforcer-ui.service || true
        
        # Reload user systemd daemon for all logged-in users
        for uid in $(loginctl list-users --no-legend | awk '{print $1}'); do
            systemctl --user --machine="${uid}@" daemon-reload 2>/dev/null || true
        done
        
        echo "USB Encryption Enforcer has been installed and started."
        echo "Configuration file: /etc/usb-enforcer/config.toml"
        echo ""
        echo "Note: Desktop notifications require the usb-enforcer-ui service to be running."
        echo "Run 'systemctl --user start usb-enforcer-ui' to start it now, or log out and back in."
        ;;
esac

#DEBHELPER#

exit 0
